# conf.d/openbox.conf
[program:openbox]
command=openbox-session
autostart=true
autorestart=true
