# conf.d/openbox.conf
[program:openbox]
command=openbox-session
autostart=true
autorestart=true
# stdout_logfile=/var/log/openbox.log
# stderr_logfile=/var/log/openbox.err.log
